# feishu-codex-bridge.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.example.feishu-codex-bridge</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/python3</string>
        <string>/ABSOLUTE/PATH/TO/OwMini/Tools/feishu_codex_bridge/bridge.py</string>
        <string>--config</string>
        <string>/ABSOLUTE/PATH/TO/OwMini/Tools/feishu_codex_bridge/bridge.toml</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/ABSOLUTE/PATH/TO/OwMini</string>

    <key>KeepAlive</key>
    <true/>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/ABSOLUTE/PATH/TO/feishu-codex-bridge.stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/ABSOLUTE/PATH/TO/feishu-codex-bridge.stderr.log</string>
</dict>
</plist>
